# docker/general/cron/crontab
* * * * * php /var/www/html/system/cron.php
